# aios-layer/deploy/systemd/aios-agent.service
[Unit]
Description=AIOS Host Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aios
Group=aios
ExecStart=/usr/local/bin/aios-agent -config /etc/aios/aios.yaml
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/aios /var/log/aios
Environment=AIOS_ENV=local

[Install]
WantedBy=multi-user.target
